# crontab
*/5 * * * * cd /workspace; /usr/bin/python ./ddns-updater.py
# EOL
