# crontab
2-57/10 * * * * cd /workspace; /usr/bin/python ./ddns-updater.py

